# poky/meta-epuck/recipes-config/epuck-rc-local-config/files/rc.local
#!/bin/sh -e

start() {
	echo "Setting USB OTG mode"
    	echo -n 1 > /sys/bus/usb/devices/3-1/bConfigurationValue
	ifup wlan0
}

stop() {
	ifdown wlan0
}

case "$1" in
	restart)
		stop
		start
		;;
	start)
		start
		;;
	stop)
		stop
		;;
	*)
		echo "Usage: /etc/init.d/rc.local {start|stop|restart}"
		exit 1
esac

exit 0
